# scripts/systemd/kmscon.service.in
[Unit]
Description=KMS System Console
Documentation=man:kmscon(1)
After=plymouth-quit-wait.service
After=systemd-user-sessions.service
After=rc-local.service

[Service]
ExecStart=@bindir@/kmscon --login -- /sbin/agetty -o '-p -- \\u' --noclear -- - $$TERM

[Install]
WantedBy=multi-user.target
